# src/services/AnaliTotTicks.service
[Unit]
Description="Get current daily percentile of stock values"

[Service]
ExecStart=/var/www/inversions/src/scripts/starter_all_stocks.sh
