# src/src/services/AnaliMainTicks.service
[Unit]
Description="Get current daily percentile of main stock values"

[Service]
ExecStart=/var/www/inversions/src/src/starter_main_stocks.sh
