# src/src/services/SortBestStocks.service
[Unit]
Description="Sort best stocks from previous day"

[Service]
ExecStart=/var/www/inversions/src/src/starter_sort_best.sh
